# deploy/lawyer.service
[Unit]
Description=Persian Legal Firm Web Application
After=network.target

[Service]
Type=simple
User=lawyer
Group=lawyer
WorkingDirectory=/srv/lawyer/current
EnvironmentFile=/etc/lawyer/lawyer.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawyer

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/lawyer
RuntimeDirectory=lawyer
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target