# deploy/lawyer.service
[Unit]
Description=Persian Legal Firm Web Application
After=network.target

[Service]
Type=simple
User=lawyer
Group=lawyer
WorkingDirectory=/home/lawyer/current
EnvironmentFile=/etc/lawyer/lawyer.env
ExecStart=/bin/bash -lc 'export HOME=/home/lawyer; export NVM_DIR=/home/lawyer/.nvm; [ -s "$NVM_DIR/nvm.sh" ] && . "$NVM_DIR/nvm.sh"; nvm use 20; exec node dist/index.js'
Environment="HOME=/home/lawyer"
Environment="NVM_DIR=/home/lawyer/.nvm"
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lawyer

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/lawyer
RuntimeDirectory=lawyer
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target